# services/stock-fetcher.service
[Unit]
Description=Taiwan Stock Data Fetcher
After=network.target

[Service]
Type=oneshot
User=YOUR_USERNAME
WorkingDirectory=/path/to/stock-strategy
ExecStart=/usr/bin/python3 /path/to/stock-strategy/fetch_taiwan_stocks.py

# 日志输出
StandardOutput=append:/var/log/stock-fetcher.log
StandardError=append:/var/log/stock-fetcher.log

# 环境变量
Environment="PYTHONUNBUFFERED=1"

# 超时设置（2小时）
TimeoutStartSec=7200

[Install]
WantedBy=multi-user.target
