# services/stock-fetcher.timer
[Unit]
Description=Taiwan Stock Data Fetcher Timer
Requires=stock-fetcher.service

[Timer]
# 每小時執行一次
OnCalendar=hourly
# 延遲隨機 0-300 秒，避免同一時間大量請求
RandomizedDelaySec=300
# 服務啟動後立即執行一次
OnBootSec=1min
# 如果錯過了執行時間，立即執行
Persistent=true

[Install]
WantedBy=timers.target
